# packages/deployment/etc/systemd/system/planv7.service
[Unit]
Description=Planv7 Application
After=network.target network-online.target

[Service]
Type=simple
Restart=always
ExecStart=/bin/sh -c 'JWT_PUBLIC_KEY=$(aws secretsmanager get-secret-value --region us-east-1 --secret-id planv7/JWT_PUBLIC_KEY) JWT_PRIVATE_KEY=$(aws secretsmanager get-secret-value --region us-east-1 --secret-id planv7/JWT_PRIVATE_KEY) /usr/bin/planv7'
