# packages/deployment/etc/systemd/system/planv7.service
[Unit]
Description=Planv7 Application
After=network.target network-online.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/planv7-server
